# supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=debug                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[unix_http_server]
file=/tmp/supervisor.sock

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:harness]
command=/harness/bootstrap 
directory=/harness
autostart=true

[program:crs]
command=/harness/crs/crs.py 0.0.0.0 56789 
directory=/harness/crs
autostart=false 

[program:irm-nova]
command=/harness/irm-nova/irm-nova.py -c /harness/cfg/irm-nova.cfg 
directory=/harness/irm-nova
autostart=false

[program:irm-shepard]
command=/harness/irm-shepard/irm-shepard.py
process_name=%(program_name)s 
directory=/harness/irm-shepard
autostart=false

[program:publisher]
command=/harness/msg-display/publisher.py  -p 6000       
directory=/harness/msg-display
autostart=false

[program:zmq_device]
command=/harness/msg-display/zmq_device.py -d FORWARDER -fb tcp://0.0.0.0:1234 -bb tcp://0.0.0.0:2345
directory=/harness/msg-display
autostart=false

[program:display]
command=/harness/msg-display/display.py -p 56787  -t 2345 -f 127.0.0.1
directory=/harness/msg-display/
autostart=false
